# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting NGINX..."

# Create required folders
mkdir -p /var/log/nginx
mkdir -p /var/lib/nginx/body
mkdir -p /var/lib/nginx/fastcgi

# Start NGINX
exec nginx 