# rootfs/etc/services.d/octo-mqtt/finish
#!/usr/bin/execlineb -P
# ==============================================================================
# Handle octo-mqtt service termination
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Import environment variables
/usr/bin/with-contenv

# Get the exit code
/usr/bin/importas -u exit_code 1

# Only bring down the container on critical failures (exit code > 128)
/usr/bin/ifelse { /usr/bin/expr ${exit_code} -gt 128 }
{ /usr/bin/ifelse { /usr/bin/expr ${exit_code} != 256 }
  { /usr/bin/foreground { /usr/bin/s6-svscanctl -t /var/run/s6/services } }
  /usr/bin/exit 0 }
/usr/bin/exit 0 