# rootfs/etc/services.d/octo-mqtt/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Take down the S6 supervision tree when octo-mqtt fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# If there are no arguments, we've terminated correctly
if { s6-test $# -eq 0 }

# If this is a non-zero exit code, something went wrong
if { s6-test ${1} -eq 256 }

# Stop the supervision tree
s6-svscanctl -t /var/run/s6/services 