# rootfs/etc/services.d/octo-mqtt/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree when octo-mqtt fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# If this is a non-zero exit code and not 256 (normal exit), bring down the container
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }
foreground { s6-svscanctl -t /var/run/s6/services }

# If this is a normal exit (256), do nothing
if { s6-test ${1} -eq 256 }
exit 0 