# rootfs/etc/services.d/octo-mqtt/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree when octo-mqtt fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# First argument is the exit code
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

# Only tell s6 to bring down the system if it isn't already doing so
if { s6-test ${1} -eq 256 }

# Stop the supervision tree
s6-svscanctl -t /var/run/s6/services 