# rootfs/etc/services.d/octo-mqtt/finish
#!/usr/bin/with-contenv bash
# ==============================================================================
# Handle octo-mqtt service termination
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Import exit code
exit_code=$1

# Log the exit
echo "Service exited with code ${exit_code}"

# Only bring down the container on critical failures (exit code > 128)
if [ "${exit_code}" -gt 128 ] && [ "${exit_code}" -ne 256 ]; then
    echo "Critical failure detected, stopping container..."
    s6-svscanctl -t /var/run/s6/services
fi

exit 0 