# rootfs/etc/services.d/octo-mqtt/run
#!/usr/bin/execlineb -P
# ==============================================================================
# Start the octo-mqtt service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Import environment variables
/usr/bin/with-contenv

# Change to the app directory and set up environment
cd /app

# Set environment variables
/usr/bin/export NODE_ENV=production
/usr/bin/export S6_KEEP_ENV=1
/usr/bin/export S6_LOGGING=1

# Start the application
/usr/bin/fdmove -c 2 1
/usr/bin/s6-setuidgid root
/usr/bin/node dist/index.js 