# rootfs/etc/services.d/octo-mqtt/run
#!/command/with-contenv
# ==============================================================================
# Start the octo-mqtt service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Change to the app directory
cd /app

# Export required environment variables
export NODE_ENV=production
export S6_KEEP_ENV=1
export S6_LOGGING=1

# Start the application
fdmove -c 2 1
exec -c
s6-setuidgid root
node dist/index.js 